# alerts/rules/imported_from_Modifying_Crontab_af202fd3-7bff-4212-a25a-fb34606cfcbe.json
{
  "rule_name": "imported_from_Modifying_Crontab_af202fd3-7bff-4212-a25a-fb34606cfcbe",
  "rule_author": "siembol",
  "rule_version": 1,
  "rule_description": "Detects suspicious modification of crontab file.",
  "source_type": "bhuio",
  "matchers": [
    {
      "is_enabled": true,
      "matcher_type": "REGEX_MATCH",
      "is_negated": false,
      "field": "original_string",
      "data": "\\QREPLACE\\E"
    }
  ],
  "tags": [
    {
      "tag_name": "sigma_id",
      "tag_value": "af202fd3-7bff-4212-a25a-fb34606cfcbe"
    },
    {
      "tag_name": "sigma_tags",
      "tag_value": "[attack.persistence, attack.t1053.003]"
    },
    {
      "tag_name": "sigma_fp",
      "tag_value": "[Legitimate modification of crontab]"
    }
  ],
  "rule_protection": {
    "max_per_hour": 30,
    "max_per_day": 100
  }
}